# rootfs/etc/service/tincd/run
#!/bin/bash

tincd -n "${TINC_NETWORK}" --config="/etc/tinc/${TINC_NETWORK}" \
    --no-detach --mlock --user=tincd --debug=${TINC_DEBUG_LEVEL}
